# playbooks/files/headnode/etc/systemd/system/netboot-persistent.target
[Unit]
Description=Netboot persistent target

[Install]
WantedBy=multi-user.target
